# imagefs/var/lib/dpkg/info/fio.list
/.
/etc
/etc/init.d
/etc/init.d/fio
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/fio.service
/usr
/usr/bin
/usr/bin/fio
/usr/bin/fio-btrace2fio
/usr/bin/fio-dedupe
/usr/bin/fio-genzipf
/usr/bin/fio2gnuplot
/usr/bin/fio_generate_plots
/usr/bin/genfio
/usr/share
/usr/share/doc
/usr/share/doc/fio
/usr/share/doc/fio/HOWTO.gz
/usr/share/doc/fio/MORAL-LICENSE
/usr/share/doc/fio/README.Debian
/usr/share/doc/fio/README.gz
/usr/share/doc/fio/REPORTING-BUGS
/usr/share/doc/fio/changelog.Debian.gz
/usr/share/doc/fio/copyright
/usr/share/doc/fio/examples
/usr/share/doc/fio/examples/1mbs_clients.fio
/usr/share/doc/fio/examples/aio-read.fio
/usr/share/doc/fio/examples/backwards-read.fio
/usr/share/doc/fio/examples/basic-verify.fio
/usr/share/doc/fio/examples/butterfly.fio
/usr/share/doc/fio/examples/cpp_null.fio
/usr/share/doc/fio/examples/cpuio.fio
/usr/share/doc/fio/examples/cross-stripe-verify.fio
/usr/share/doc/fio/examples/dev-dax.fio
/usr/share/doc/fio/examples/disk-zone-profile.fio
/usr/share/doc/fio/examples/e4defrag.fio
/usr/share/doc/fio/examples/e4defrag2.fio
/usr/share/doc/fio/examples/enospc-pressure.fio
/usr/share/doc/fio/examples/falloc.fio
/usr/share/doc/fio/examples/filecreate-ioengine.fio
/usr/share/doc/fio/examples/fio-rand-RW.fio
/usr/share/doc/fio/examples/fio-rand-read.fio
/usr/share/doc/fio/examples/fio-rand-write.fio
/usr/share/doc/fio/examples/fio-seq-RW.fio
/usr/share/doc/fio/examples/fio-seq-read.fio
/usr/share/doc/fio/examples/fio-seq-write.fio
/usr/share/doc/fio/examples/fixed-rate-submission.fio
/usr/share/doc/fio/examples/flow.fio
/usr/share/doc/fio/examples/fsx.fio
/usr/share/doc/fio/examples/ftruncate.fio
/usr/share/doc/fio/examples/gfapi.fio
/usr/share/doc/fio/examples/gpudirect-rdmaio-client.fio
/usr/share/doc/fio/examples/gpudirect-rdmaio-server.fio
/usr/share/doc/fio/examples/http-s3.fio
/usr/share/doc/fio/examples/http-swift.fio
/usr/share/doc/fio/examples/http-webdav.fio
/usr/share/doc/fio/examples/ime.fio
/usr/share/doc/fio/examples/iometer-file-access-server.fio
/usr/share/doc/fio/examples/jesd219.fio
/usr/share/doc/fio/examples/latency-profile.fio
/usr/share/doc/fio/examples/libhdfs.fio
/usr/share/doc/fio/examples/libpmem.fio
/usr/share/doc/fio/examples/mtd.fio
/usr/share/doc/fio/examples/netio.fio
/usr/share/doc/fio/examples/netio_multicast.fio
/usr/share/doc/fio/examples/null.fio
/usr/share/doc/fio/examples/numa.fio
/usr/share/doc/fio/examples/pmemblk.fio
/usr/share/doc/fio/examples/poisson-rate-submission.fio
/usr/share/doc/fio/examples/rados.fio
/usr/share/doc/fio/examples/rand-zones.fio
/usr/share/doc/fio/examples/rbd.fio
/usr/share/doc/fio/examples/rdmaio-client.fio
/usr/share/doc/fio/examples/rdmaio-server.fio
/usr/share/doc/fio/examples/ssd-steadystate.fio
/usr/share/doc/fio/examples/ssd-test.fio
/usr/share/doc/fio/examples/steadystate.fio
/usr/share/doc/fio/examples/surface-scan.fio
/usr/share/doc/fio/examples/tiobench-example.fio
/usr/share/doc/fio/examples/waitfor.fio
/usr/share/doc/fio/examples/zipf.fio
/usr/share/doc/fio/fiologparser.py.gz
/usr/share/doc/fio/fiologparser_hist.py.1.gz
/usr/share/doc/fio/fiologparser_hist.py.gz
/usr/share/doc/fio/half-bins.py
/usr/share/doc-base
/usr/share/doc-base/fio
/usr/share/fio
/usr/share/fio/graph2D.gpm
/usr/share/fio/graph3D.gpm
/usr/share/fio/math.gpm
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/fio.1.gz
/usr/share/man/man1/fio2gnuplot.1.gz
/usr/share/man/man1/fio_generate_plots.1.gz
