# imagefs/var/lib/dpkg/info/fwupd.list
/.
/etc
/etc/fwupd
/etc/fwupd/daemon.conf
/etc/fwupd/redfish.conf
/etc/fwupd/remotes.d
/etc/fwupd/remotes.d/dell-esrt.conf
/etc/fwupd/remotes.d/lvfs-testing.conf
/etc/fwupd/remotes.d/lvfs.conf
/etc/fwupd/remotes.d/vendor-directory.conf
/etc/fwupd/remotes.d/vendor.conf
/etc/fwupd/thunderbolt.conf
/etc/fwupd/uefi.conf
/etc/fwupd/upower.conf
/etc/pki
/etc/pki/fwupd
/etc/pki/fwupd/GPG-KEY-Hughski-Limited
/etc/pki/fwupd/GPG-KEY-Linux-Foundation-Firmware
/etc/pki/fwupd/GPG-KEY-Linux-Vendor-Firmware-Service
/etc/pki/fwupd/LVFS-CA.pem
/etc/pki/fwupd-metadata
/etc/pki/fwupd-metadata/GPG-KEY-Linux-Foundation-Metadata
/etc/pki/fwupd-metadata/GPG-KEY-Linux-Vendor-Firmware-Service
/etc/pki/fwupd-metadata/LVFS-CA.pem
/etc/update-motd.d
/etc/update-motd.d/85-fwupd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/fwupd-offline-update.service
/lib/systemd/system/fwupd-refresh.service
/lib/systemd/system/fwupd-refresh.timer
/lib/systemd/system/fwupd.service
/lib/systemd/system/system-update.target.wants
/lib/systemd/system-preset
/lib/systemd/system-preset/fwupd-refresh.preset
/lib/systemd/system-shutdown
/lib/systemd/system-shutdown/fwupd.shutdown
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/90-fwupd-devices.rules
/usr
/usr/bin
/usr/bin/dfu-tool
/usr/bin/fwupdagent
/usr/bin/fwupdate
/usr/bin/fwupdmgr
/usr/bin/fwupdtool
/usr/bin/fwupdtpmevlog
/usr/lib
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_altos.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_amt.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_ata.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_colorhug.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_coreboot.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_csr.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_dell.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_dell_dock.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_dell_esrt.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_dfu.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_ebitdo.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_emmc.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_fastboot.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_fresco_pd.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_jabra.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_logind.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_logitech_hidpp.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_nitrokey.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_nvme.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_optionrom.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_redfish.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_rts54hid.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_rts54hub.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_solokey.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_steelseries.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_superio.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_synaptics_cxaudio.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_synaptics_mst.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_synaptics_prometheus.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_synaptics_rmi.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_thelio_io.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_thunderbolt.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_thunderbolt_power.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_tpm.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_tpm_eventlog.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_uefi.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_uefi_recovery.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_upower.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_vli.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_wacom_raw.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-3/libfu_plugin_wacom_usb.so
/usr/libexec
/usr/libexec/fwupd
/usr/libexec/fwupd/efi
/usr/libexec/fwupd/efi/fwupdx64.efi
/usr/libexec/fwupd/fwupd
/usr/libexec/fwupd/fwupdoffline
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/fwupdagent
/usr/share/bash-completion/completions/fwupdmgr
/usr/share/bash-completion/completions/fwupdtool
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.fwupd.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.fwupd.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/org.freedesktop.fwupd.conf
/usr/share/doc
/usr/share/doc/fwupd
/usr/share/doc/fwupd/copyright
/usr/share/fish
/usr/share/fish/vendor_completions.d
/usr/share/fish/vendor_completions.d/fwupdmgr.fish
/usr/share/fwupd
/usr/share/fwupd/add_capsule_header.py
/usr/share/fwupd/firmware_packager.py
/usr/share/fwupd/install_dell_bios_exe.py
/usr/share/fwupd/metainfo
/usr/share/fwupd/metainfo/org.freedesktop.fwupd.remotes.lvfs-testing.metainfo.xml
/usr/share/fwupd/metainfo/org.freedesktop.fwupd.remotes.lvfs.metainfo.xml
/usr/share/fwupd/quirks.d
/usr/share/fwupd/quirks.d/altos.quirk
/usr/share/fwupd/quirks.d/ata.quirk
/usr/share/fwupd/quirks.d/colorhug.quirk
/usr/share/fwupd/quirks.d/coreboot.quirk
/usr/share/fwupd/quirks.d/csr-aiaiai.quirk
/usr/share/fwupd/quirks.d/dell-dock.quirk
/usr/share/fwupd/quirks.d/dell.quirk
/usr/share/fwupd/quirks.d/dfu.quirk
/usr/share/fwupd/quirks.d/ebitdo.quirk
/usr/share/fwupd/quirks.d/emmc.quirk
/usr/share/fwupd/quirks.d/fastboot.quirk
/usr/share/fwupd/quirks.d/fresco-pd.quirk
/usr/share/fwupd/quirks.d/jabra.quirk
/usr/share/fwupd/quirks.d/logitech-hidpp.quirk
/usr/share/fwupd/quirks.d/nitrokey.quirk
/usr/share/fwupd/quirks.d/nvme.quirk
/usr/share/fwupd/quirks.d/optionrom.quirk
/usr/share/fwupd/quirks.d/rts54hid.quirk
/usr/share/fwupd/quirks.d/rts54hub.quirk
/usr/share/fwupd/quirks.d/solokey.quirk
/usr/share/fwupd/quirks.d/steelseries.quirk
/usr/share/fwupd/quirks.d/superio.quirk
/usr/share/fwupd/quirks.d/synaptics-cxaudio.quirk
/usr/share/fwupd/quirks.d/synaptics-mst.quirk
/usr/share/fwupd/quirks.d/synaptics-prometheus.quirk
/usr/share/fwupd/quirks.d/synaptics-rmi.quirk
/usr/share/fwupd/quirks.d/thelio-io.quirk
/usr/share/fwupd/quirks.d/tpm.quirk
/usr/share/fwupd/quirks.d/uefi-recovery.quirk
/usr/share/fwupd/quirks.d/uefi.quirk
/usr/share/fwupd/quirks.d/vli-pd.quirk
/usr/share/fwupd/quirks.d/vli-usbhub-lenovo.quirk
/usr/share/fwupd/quirks.d/vli-usbhub.quirk
/usr/share/fwupd/quirks.d/wacom-raw.quirk
/usr/share/fwupd/quirks.d/wacom-usb.quirk
/usr/share/fwupd/remotes.d
/usr/share/fwupd/remotes.d/dell-esrt
/usr/share/fwupd/remotes.d/dell-esrt/metadata.xml
/usr/share/fwupd/remotes.d/vendor
/usr/share/fwupd/remotes.d/vendor/firmware
/usr/share/fwupd/remotes.d/vendor/firmware/README.md
/usr/share/fwupd/simple_client.py
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/scalable
/usr/share/icons/hicolor/scalable/apps
/usr/share/icons/hicolor/scalable/apps/org.freedesktop.fwupd.svg
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/fwupd
/usr/share/locale
/usr/share/locale/ca
/usr/share/locale/ca/LC_IMAGES
/usr/share/locale/ca/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/ca/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/ca/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/ca/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/ca/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/ca/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/ca/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/ca/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/cs
/usr/share/locale/cs/LC_IMAGES
/usr/share/locale/cs/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/cs/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/cs/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/cs/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/cs/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/cs/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/cs/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/cs/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/da
/usr/share/locale/da/LC_IMAGES
/usr/share/locale/da/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/da/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/da/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/da/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/da/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/da/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/da/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/da/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/de
/usr/share/locale/de/LC_IMAGES
/usr/share/locale/de/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/de/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/de/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/de/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/de/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/de/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/de/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/de/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/en
/usr/share/locale/en/LC_IMAGES
/usr/share/locale/en/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/en/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/en/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/en/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/en/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/en/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/en/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/en/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/fi
/usr/share/locale/fi/LC_IMAGES
/usr/share/locale/fi/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/fi/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/fi/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/fi/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/fi/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/fi/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/fi/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/fi/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/fur
/usr/share/locale/fur/LC_IMAGES
/usr/share/locale/fur/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/fur/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/fur/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/fur/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/fur/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/fur/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/fur/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/fur/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/hr
/usr/share/locale/hr/LC_IMAGES
/usr/share/locale/hr/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/hr/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/hr/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/hr/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/hr/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/hr/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/hr/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/hr/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/hu
/usr/share/locale/hu/LC_IMAGES
/usr/share/locale/hu/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/hu/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/hu/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/hu/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/hu/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/hu/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/hu/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/hu/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/id
/usr/share/locale/id/LC_IMAGES
/usr/share/locale/id/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/id/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/id/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/id/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/id/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/id/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/id/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/id/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/it
/usr/share/locale/it/LC_IMAGES
/usr/share/locale/it/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/it/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/it/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/it/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/it/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/it/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/it/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/it/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/ko
/usr/share/locale/ko/LC_IMAGES
/usr/share/locale/ko/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/ko/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/ko/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/ko/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/ko/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/ko/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/ko/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/ko/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/lt
/usr/share/locale/lt/LC_IMAGES
/usr/share/locale/lt/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/lt/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/lt/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/lt/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/lt/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/lt/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/lt/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/lt/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/pl
/usr/share/locale/pl/LC_IMAGES
/usr/share/locale/pl/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/pl/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/pl/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/pl/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/pl/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/pl/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/pl/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/pl/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_IMAGES
/usr/share/locale/pt_BR/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/pt_BR/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/pt_BR/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/pt_BR/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/pt_BR/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/pt_BR/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/pt_BR/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/pt_BR/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/ru
/usr/share/locale/ru/LC_IMAGES
/usr/share/locale/ru/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/ru/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/ru/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/ru/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/ru/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/ru/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/ru/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/ru/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/sr
/usr/share/locale/sr/LC_IMAGES
/usr/share/locale/sr/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/sr/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/sr/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/sr/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/sr/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/sr/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/sr/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/sr/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/sv
/usr/share/locale/sv/LC_IMAGES
/usr/share/locale/sv/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/sv/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/sv/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/sv/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/sv/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/sv/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/sv/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/sv/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/tr
/usr/share/locale/tr/LC_IMAGES
/usr/share/locale/tr/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/tr/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/tr/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/tr/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/tr/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/tr/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/tr/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/tr/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/uk
/usr/share/locale/uk/LC_IMAGES
/usr/share/locale/uk/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/uk/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/uk/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/uk/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/uk/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/uk/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/uk/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/uk/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/zh_CN
/usr/share/locale/zh_CN/LC_IMAGES
/usr/share/locale/zh_CN/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/zh_CN/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/zh_CN/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/zh_CN/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/zh_CN/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/zh_CN/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/zh_CN/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/zh_CN/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/locale/zh_TW
/usr/share/locale/zh_TW/LC_IMAGES
/usr/share/locale/zh_TW/LC_IMAGES/fwupd-1024-768.bmp.gz
/usr/share/locale/zh_TW/LC_IMAGES/fwupd-1920-1080.bmp.gz
/usr/share/locale/zh_TW/LC_IMAGES/fwupd-3840-2160.bmp.gz
/usr/share/locale/zh_TW/LC_IMAGES/fwupd-5120-2880.bmp.gz
/usr/share/locale/zh_TW/LC_IMAGES/fwupd-5688-3200.bmp.gz
/usr/share/locale/zh_TW/LC_IMAGES/fwupd-640-480.bmp.gz
/usr/share/locale/zh_TW/LC_IMAGES/fwupd-7680-4320.bmp.gz
/usr/share/locale/zh_TW/LC_IMAGES/fwupd-800-600.bmp.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/dfu-tool.1.gz
/usr/share/man/man1/fwupdagent.1.gz
/usr/share/man/man1/fwupdate.1.gz
/usr/share/man/man1/fwupdmgr.1.gz
/usr/share/man/man1/fwupdtool.1.gz
/usr/share/man/man1/fwupdtpmevlog.1.gz
/usr/share/metainfo
/usr/share/metainfo/org.freedesktop.fwupd.metainfo.xml
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.fwupd.policy
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/org.freedesktop.fwupd.rules
/var
/var/cache
/var/cache/app-info
/var/cache/app-info/xmls
/var/lib
/var/lib/fwupd
/var/lib/fwupd/builder
/var/lib/fwupd/builder/README.md
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/fwupd.pkla
/lib/systemd/system/system-update.target.wants/fwupd-offline-update.service
/usr/share/doc/fwupd/changelog.Debian.gz
