# imagefs/var/lib/dpkg/info/multipath-tools.conffiles
/etc/init.d/multipath-tools
/etc/init.d/multipath-tools-boot
